# usbip/usbip-Client.old
#!/bin/bash

### BEGIN INIT INFO
# Provides:          usbip-Client
# Required-Start:    $local_fs $network $named $time $syslog
# Required-Stop:     $local_fs $network $named $time $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Description:       P5Software USB IP Client
### END INIT INFO

# ------------------------------------------------------------
# Configuration
# ------------------------------------------------------------
P5SoftwareHome=/etc/P5Software/Linux-Tools
scripthome=$P5SoftwareHome/usbip

# Read the configuration file
source $scripthome/usbip-Client.conf

SCRIPT="usbip attach -r $serverFQDN -b $busID"
RUNAS=root

PIDFILE=/var/run/usbip-Client.pid
LOGFILE=/var/log/usbip-Client.log

start() {
if [ -f /var/run/$PIDNAME ] && kill -0 $(cat /var/run/$PIDNAME); then
echo $(date +%Y-%m-%d\ %H:%M:%S) 'Service already running' >> $LOGFILE 2>&1
return 1
fi

#Set the environment appropriately
modprobe vhci-hcd

echo $(date +%Y-%m-%d\ %H:%M:%S) 'Starting service…' >> $LOGFILE 2>&1
local CMD="$SCRIPT >> $LOGFILE 2>&1"
su -c "$CMD" $RUNAS > $PIDFILE
echo $(date +%Y-%m-%d\ %H:%M:%S) 'Service started' >> $LOGFILE 2>&1
}

stop() {
if [ ! -f $PIDFILE ] || ! kill -0 $(cat $PIDFILE); then
echo $(date +%Y-%m-%d\ %H:%M:%S) 'Service not running' >> $LOGFILE 2>&1
return 1
fi
echo $(date +%Y-%m-%d\ %H:%M:%S) 'Stopping service…' >> $LOGFILE 2>&1

#Unbind the device
local CMD="usbip detach -p 0 >> $LOGFILE 2>&1"
su -c "$CMD" $RUNAS >> $LOGFILE 2>&1
kill -15 $(cat "$PIDFILE") && rm -f $PIDFILE
echo $(date +%Y-%m-%d\ %H:%M:%S) 'Service stopped' >> $LOGFILE 2>&1
}


case "$1" in
start)
start
;;
stop)
stop
;;
uninstall)
uninstall
;;
retart)
stop
start
;;
*)
echo "Usage: $0 {start|stop|restart|uninstall}"
esac